# monthly-reboot/monthly-reboot.service
[Unit]
Description=Monthly System Reboot
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStartPre=sleep 60
ExecStart=reboot

[Install]
WantedBy=multi-user.target
